# throttle_stop/etc/systemd/system/undervolt.service
[Unit]
Description=Undervolt Service

[Service]
ExecStart=/usr/local/bin/undervolt.sh

[Install]
WantedBy=multi-user.target
